# django/docker/cron/crontab
*/3 * * * * /bin/bash /home/flash/cronjobs/delete_django_files.sh
# Don't remove the empty line at the end of this file. It is required to run the cron job
